# backup/systemd/kitty-theme-reload.service
[Unit]
Description=Reload kitty theme when caelestia scheme changes

[Service]
Type=oneshot
ExecStart=%h/.config/caelestia/scripts/kitty-theme-reload.sh
