# backup/systemd/kitty-theme-reload.path
[Unit]
Description=Watch caelestia kitty theme file for changes

[Path]
PathModified=%h/.local/state/caelestia/theme/kitty-theme.conf

[Install]
WantedBy=default.target
